# src/argus_alert/supervisor.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp

[program:driver]
command=python3 driver.py
directory=/home/argus/argus-alert/argus_alert/core/inspect/
numprocs=1
priority=1
autostart=true

[program:exector]
command=python3 executor.py
directory=/home/argus/argus-alert/argus_alert/core/inspect/
numprocs=1
autostart=true

[program:notifier]
command=python3 handler.py
priority=1
directory=/home/argus/argus-alert/argus_alert/core/notice/
numprocs=1
autostart=true

